# deploy/doors.service
[Unit]
Description=DOORS Game Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/doors
Environment=NODE_ENV=production
Environment=PORT=3400
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
